# root/etc/systemd/system/suspend_low_bat.timer
[Unit]
Description=Check battery every minute, suspending if below 15%

[Timer]
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
